# services/navidrome.third-branches.net/systemd/navidrome-tunnel.container
[Unit]
Description=Cloudflare tunnel to navidrome.third-branches.net
Wants=navidrome-app.container
After=navidrome-app.container

[Service]
Restart=always
SyslogLevel=info
TimeoutStartSec=900

[Container]
Label="navidrome-tunnel"
Pod="navidrome.pod"
Image=docker.io/cloudflare/cloudflared:latest
Environment=TUNNEL_TOKEN= # TODO: Set tunnel token generated by opentofu.
Exec=tunnel run

[Install]
WantedBy=default.target
